# deploy/delta-neutral-monitor-webapp.service
[Unit]
Description=Delta Neutral Monitor Next.js webapp
After=network-online.target delta-neutral-monitor-backend.service
Wants=network-online.target
Requires=delta-neutral-monitor-backend.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/projects/delta-neutral-monitor/webapp
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3000
Environment=MONITOR_API_BASE_URL=http://127.0.0.1:8080
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
